# services/lifecycle_manager/systemd/valence-lifecycle-manager.service
[Unit]
Description=Valence Lifecycle Manager Service
After=network.target postgresql.service

[Service]
Type=simple
User=valence
Group=valence
WorkingDirectory=/opt/valence/lifecycle-manager
ExecStart=/opt/valence/lifecycle-manager/valence-lifecycle-manager
Restart=always
RestartSec=10

# Environment
Environment="RUST_LOG=info"
EnvironmentFile=/etc/valence/lifecycle-manager.env

# Security
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/var/lib/valence

# Resource limits
LimitNOFILE=65536
MemoryLimit=2G
CPUQuota=200%

[Install]
WantedBy=multi-user.target